# duckdns/samples/osx/crontab
*/5 * * * * /Users/administrator/duckdns/duck.sh