# duckdns/samples/linux_cron/crontab
*/5 * * * * /home/douglas/duckdns/duck.sh